# systemd/phantom-monitor.service
[Unit]
Description=Phantom Teams Monitor
Documentation=https://github.com/NinjaTech-AI/phantom
After=network.target supervisor.service
Wants=supervisor.service

[Service]
Type=simple
User=root
WorkingDirectory=/workspace
ExecStartPre=/bin/sleep 10
ExecStart=/usr/local/bin/python -m phantom.teams_monitor --interval 15
Restart=always
RestartSec=30s

Environment=HOME=/root
Environment=DISPLAY=:99
Environment=PYTHONPATH=/workspace
StandardOutput=append:/workspace/logs/phantom-teams-monitor.log
StandardError=append:/workspace/logs/phantom-teams-monitor.log
SyslogIdentifier=phantom-teams-monitor

[Install]
WantedBy=multi-user.target
